# supervisor/supervisor.ini
[group:easycrawl]
programs=easycrawl.beat,easycrawl.worker

[program:easycrawl.beat]
command=%(ENV_HOME)s/venv/easycrawl/bin/celery -A task  beat -l info
priority=999
autostart=true
startsecs=10
exitcodes=0,2
directory=%(ENV_HOME)s/project/easycrawl
user=yudan
redirect_stderr=true
stdout_logfile=%(ENV_HOME)s/logs/easycrawl.beat.log
stdout_logfile_maxbytes=250MB
stdout_logfile_backups=10
stderr_logfile=%(ENV_HOME)s/logs/easycrawl.beat.err
stderr_logfile_maxbytes=250MB
stderr_logfile_backups=10
environment=PYTHONPATH='%(ENV_HOME)s/project/easycrawl';LANG=en_US.utf8

[program:easycrawl.worker]
command=%(ENV_HOME)s/venv/easycrawl/bin/celery -A task worker -l info -n easycrawl.worker
directory=%(ENV_HOME)s/project/easycrawl
autostart=true
user=yudan
startsecs=10
stdout_logfile=%(ENV_HOME)s/logs/easycrawl.worker.log
stdout_logfile_maxbytes=250MB
stdout_logfile_backups=10
stderr_logfile=%(ENV_HOME)s/logs/easycrawl.worker.err
stderr_logfile_maxbytes=250MB
stderr_logfile_backups=10
environment=PYTHONPATH='%(ENV_HOME)s/project/easycrawl';LANG=en_US.utf8


